# supporting-files/wpsd-running-tasks.timer
[Unit]
Description=Run WPSD Tasks at random intervals (max 1 hour apart)

[Timer]
# Wait 1 minute after the service finishes
OnUnitInactiveSec=1m
# Then add a random delay of up to 59 minutes
RandomizedDelaySec=59m

[Install]
WantedBy=timers.target
